# supervisor/config/supervisord.ini
[supervisord]
logfile=/var/log/supervisor/supervisord.log
pidfile=/run/supervisord.pid

[program:php]
command=/usr/local/sbin/php-fpm

[program:cron]
command=/usr/sbin/crond -L /var/log/cron/cron.log -f